# backup.service
[Unit]
Description=Contestant home folder backup service

[Service]
User=ansible
ExecStart=/bin/bash /home/ansible/backup.sh
