# backup.service
[Unit]
Description=Contestant home folder backup service

[Service]
ExecStart=/bin/bash ~/backup.sh