# backup.timer
[Unit]
Description=Contestant home folder backup timer

[Timer]
OnUnitActiveSec=5m
Unit=backup.service

[Install]
WantedBy=multi-user.target
