# isolate-install-configure/backup/backup.timer
[Unit]
Description=Contestant home folder backup timer

[Timer]
OnUnitActiveSec=10m
OnBootSec=10s
Unit=backup.service

[Install]
WantedBy=timers.target
